# scripts/delphinet-deploy.service
[Unit]
Description=DelphiNet 6 — auto-deploy on new git tags
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# These match the defaults produced by scripts/server-install.sh.
# If you ran the install script with custom RUN_USER or INSTALL_DIR overrides,
# the install script will sed-replace these lines automatically.
User=delphinet
WorkingDirectory=/opt/delphinet6
ExecStart=/bin/bash /opt/delphinet6/scripts/watch-deploy.sh

Restart=on-failure
RestartSec=10

# Override poll interval here or via:
#   sudo systemctl edit delphinet-deploy
#   [Service]
#   Environment="POLL_INTERVAL=60"
Environment="POLL_INTERVAL=30"
Environment="REPO_DIR=/opt/delphinet6"

# View logs with: journalctl -u delphinet-deploy -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=delphinet-deploy

[Install]
WantedBy=multi-user.target
